# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=$HOME/.config/polybar/launch.sh
Hidden=false
X-MATE-Autostart-enabled=true
Name=Polybar
