# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=bash $HOME/.config/polybar/launch.sh
Hidden=false
Name=Polybar
